# rootfs/etc/service.d/event-client/run
#!/usr/bin/execlineb -P

cd /app
node client.js
