# rootfs/etc/services.d/event-client/run
#!/usr/bin/with-contenv sh

echo $TUTUM_USER
echo $TUTUM_APIKEY

node /app/src/client.js
